# parse.cron
0 0 * * * python parse.py
0 6 * * * python parse.py
0 12 * * * python parse.py
0 18 * * * python parse.py
