# parse.cron
00 00 * * * python parse.py
00 06 * * * python parse.py
00 12 * * * python parse.py
00 18 * * * python parse.py
